# devops/weight_billing_monitor.service
[Unit]
Description=Weight and Billing Monitoring Service
After=network.target

[Service]
Type=simple
Restart=always
ExecStart=/automation/health_monitor.sh

[Install]
WantedBy=multi-user.target